# systemd/buddyhub.service
# Claude Buddy hub — systemd --user unit (example).
#
# Install:
#   mkdir -p ~/.config/systemd/user
#   cp systemd/buddyhub.service ~/.config/systemd/user/
#   # edit WorkingDirectory / ExecStart paths and --owner below
#   systemctl --user daemon-reload
#   systemctl --user enable --now buddyhub.service
#   loginctl enable-linger "$USER"     # keep it running across logout
#
# Adjust the two %h-relative paths if you cloned buddy-bridge elsewhere.

[Unit]
Description=Claude Buddy hub (buddy-bridge)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/git/buddy-bridge
# -u keeps stdout unbuffered so relay/device events show live in journalctl.
ExecStart=/usr/bin/python3 -u %h/git/buddy-bridge/buddyhub.py --port 8787 --transport relay --owner you
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
