# systemd/buddy-tunnel.service
# Claude Buddy SSH reverse tunnel — systemd --user unit (example).
#
# Exposes the hub's :8787 on a remote machine's localhost:8787 so that machine's
# Claude Code hooks can reach the hub without opening any inbound ports. Only
# needed if you run sessions on a second box; skip it for a single-machine setup.
#
# Install:
#   cp systemd/buddy-tunnel.service ~/.config/systemd/user/
#   # replace `remote-host` below with your SSH host/alias (the other machine)
#   systemctl --user daemon-reload
#   systemctl --user enable --now buddy-tunnel.service
#
# Requires key-based SSH auth (no password prompt) to that host.

[Unit]
Description=Claude Buddy SSH reverse tunnel (buddy-bridge)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/ssh -N -o ServerAliveInterval=30 -o ServerAliveCountMax=3 -o ExitOnForwardFailure=yes -R 8787:localhost:8787 remote-host
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
